# .config/systemd/wallust-theme.path
[Unit]
Description=Checa se o config.ini do waypaper mudou e atualiza os temas

[Path]
PathChanged=%h/.config/waypaper/config.ini
Unit=wallust-theme.service

[Install]
WantedBy=default.target
